# roku.service
[Unit]
Description=RokuAlexa Lamba Skill Container as a service
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker run --net=host --name rokualexa jaknoll/rokualexa:1.0.0
ExecStop=/usr/bin/docker rm -f rokualexa 
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
